# filesdir/.config/systemd/user/r2e.service
[Unit]
Description=rss2email service

[Service]
Type=oneshot
ExecStart=/usr/bin/r2e run
